# mtls_server/assets/api_server
#!/bin/sh /etc/rc.common

T=99
STOP=01

start() {
	chmod +x /bin/easynet/api_server
	/bin/easynet/api_server --dashboard-path $(pwd)/dashboard --server-ip 0.0.0.0 &
}

stop() {
	killall api_server
}
